# crontab/syetem_update.sh
#!/bin/bash
sudo apt-get clean
sudo apt-get update -qq
sudo apt-get upgrade -qqy
sudo apt-get autoclean -qq
